# advise-backend-v2.service
[Unit]
Description=Advise Backend Service V2

[Service]
ExecStart=/home/zonainovasi/backend_services/advise-backend-v2/start.sh
